# crontab
# M  H  DM M  DW|
#---------------+
  0  0  *  *  *   /bin/bash /opt/local/bin/daily-tarsnap.sh >/dev/null 2>&1